# service/kafka_shell.ini
[program:kafka_shell]
directory=/data/kafkamonitor/media
command=/usr/bin/python kafkainfo_monitor.py
numprocs=1
autostart=true
autorestart=true